# unit_files/service6.service.d/custom.conf
[Unit]
Before=nginx.service

[Service]
EnvironmentFile=/root/envfile.conf
Environment='TESTE=Esse eh o valor da variavel de ambiente TESTE'
